# rootfs/etc/s6-overlay/s6-rc.d/tracker/run
#!/command/with-contenv sh
exec crond -f -l 2
